# deploy/jibe.service
[Unit]
Description=Jibe daemon — Android/Linux LAN bridge
Documentation=https://github.com/gobbledyglomp/jibe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/jibe
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jibe

[Install]
WantedBy=default.target
